# uxplay_ble_service/uxplay.timer
[Unit]
Description=Timer that initially starts uxplay.service

[Timer]
OnStartupSec=30
Unit=uxplay.service

[Install]
WantedBy=default.target
